# gnome-mqtt-tray.service
[Unit]
Description=Gnome Mqtt Tray Service

[Service]
Type=simple
ExecStart=/usr/bin/gnome-mqtt-tray.py

[Install]
WantedBy=multi-user.target
